# install/systemd/jabali-backup-retention.service
[Unit]
Description=Jabali backup retention sweep (M30 — restic forget + prune)
Documentation=https://git.linux-hosting.co.il/shukivaknin/jabali2/src/branch/main/docs/adr/0075-backup-restore-restic.md
After=jabali-panel.service network-online.target
# No Requires=jabali-panel.service. The CLI reads server_settings to pick
# up keep-daily/weekly/monthly knobs but the panel may be transiently down
# during boot; install.sh enables this timer in the install foundation
# step before the panel unit lands. A failed run is tolerated — the timer
# fires again at 04:30 the next day. `Restart=no` matches the M22 reaper.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jabali backup retention apply
StandardOutput=journal
StandardError=journal

# Restic needs a writable cache dir. ProtectHome=read-only blocks
# /root/.cache, so point RESTIC_CACHE_DIR at a path under
# ReadWritePaths instead. HOME stays pinned so any restic codepath
# that falls back to $HOME also lands somewhere sane.
Environment=HOME=/root
Environment=RESTIC_CACHE_DIR=/var/lib/jabali-backups/.cache/restic

# Restic prune walks every pack file in /var/lib/jabali-backups/repo/. On
# a multi-GB repo this can hit several minutes of CPU + disk I/O; the
# 04:30 schedule keeps it off-peak. ProtectHome=read-only because the
# retention sweep never touches /home; it only purges restic packs.
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/jabali-backups
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Runs as root (the default — no User= directive). Rationale: the
# agent invokes restic as root for every backup (it has to traverse
# /home/<user>/* which only root can read), so restic creates
# 0600 root-owned packs in /var/lib/jabali-backups/repo/. Retention
# must read+write those same packs to forget+prune; running as root
# matches the agent's identity and avoids a chown -R after every
# backup. The systemd hardening above (PrivateTmp, ProtectSystem,
# ReadWritePaths=/var/lib/jabali-backups) keeps the blast radius
# narrow.

Restart=no
