# install/systemd/jabali-hydra.service
[Unit]
Description=Ory Hydra OAuth 2 / OIDC Provider
Documentation=https://www.ory.sh/docs/hydra
# Kratos must be reachable for the first OIDC flow to complete
# (login/consent handlers in panel-api delegate to Kratos session
# validation). No MariaDB dependency — Hydra's state lives in a
# SQLite file under StateDirectory, see ADR-0036 amendment.
After=network-online.target jabali-kratos.service
Wants=network-online.target
PartOf=jabali.slice

[Service]
Type=simple
User=jabali
Group=jabali
Slice=jabali.slice

# StateDirectory ensures /var/lib/jabali-hydra/ exists with the correct
# owner on every boot. First boot creates it; subsequent boots no-op.
# systemd also sets the directory permissions per StateDirectoryMode.
StateDirectory=jabali-hydra
StateDirectoryMode=0750

# Hydra runs its own migrations on every boot; idempotent + safe to
# re-run. SQLite migrations are upstream-clean (no CAST AS JSON
# landmines that MariaDB rejects).
ExecStart=/usr/local/bin/hydra serve all -c /etc/jabali-panel/hydra.yml

Restart=on-failure
RestartSec=5s

# Security hardening mirrors jabali-kratos.service. ReadWritePaths
# covers the config dir (future live-reload of signing keys), the
# SQLite state dir, and /run/jabali for any runtime scratch.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/jabali-panel /var/lib/jabali-hydra /run/jabali

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
