# install/systemd/jabali-backup-retention.timer
[Unit]
Description=Daily restic forget+prune for Jabali backup repo
Documentation=https://git.linux-hosting.co.il/shukivaknin/jabali2/src/branch/main/docs/adr/0075-backup-restore-restic.md

[Timer]
# 04:30 host-local time. Off-peak; long enough after 03:00 maldet
# scans + 04:00 quarantine purge to avoid I/O contention.
OnCalendar=*-*-* 04:30:00
# Fire missed runs after a host wakes from sleep / suspend.
Persistent=true
Unit=jabali-backup-retention.service

[Install]
WantedBy=timers.target
