# install/systemd/jabali-migration-secrets-reap.service
[Unit]
Description=Wipe migration-secrets env files for terminal-state migration_jobs (M35 ADR-0094)
Documentation=https://git.linux-hosting.co.il/shukivaknin/jabali2/src/branch/main/docs/adr/0094-migration-importers.md
After=jabali-panel.service mariadb.service
# No Requires= — daily cadence + Persistent=true means a missed run on
# a transiently-unhealthy panel just slides to the next tick; same
# resilience pattern as jabali-sso-reaper.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/jabali migrate reap-secrets
StandardOutput=journal
StandardError=journal
Restart=no

# Hardening: reaper only needs read on migration_jobs (already wired
# through the panel DB cred) + unlink on /etc/jabali-panel/migration-
# secrets/*.env. ProtectSystem + ReadWritePaths confines the surface.
ProtectSystem=strict
ReadWritePaths=/etc/jabali-panel/migration-secrets
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
