# install/systemd/jabali-filebrowser.service
[Unit]
Description=Jabali filebrowser service — reverse-proxy target for file management
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# -b /files is pinned on the command line because filebrowser persists
# the baseurl into its SQLite DB on first run and thereafter ignores
# config.json changes. CLI flags override the DB-saved value, so this
# is the only reliable way to keep nginx's /files prefix stable.
ExecStart=/usr/local/bin/filebrowser -c /etc/jabali-panel/filebrowser/config.json -b /files
Restart=on-failure
RestartSec=5s

# Runs as dedicated filebrowser system user.
# Reconciler adds each jabali user's group via usermod -aG <user> filebrowser,
# allowing filebrowser to read/write per-user home directories (mode 0750).
# See ADR-0027 §7 (Open Question #1 — resolved to option b).
User=filebrowser
Group=filebrowser

# Use RuntimeDirectory to auto-create /run/jabali-filebrowser with mode 0755
RuntimeDirectory=jabali-filebrowser

# Use StateDirectory to auto-create /var/lib/jabali-filebrowser for the filebrowser DB
StateDirectory=jabali-filebrowser

# Security hardening: filebrowser accesses /home for user files (via group membership)
# and /run/jabali-filebrowser for the UDS socket. Does NOT listen on TCP.
# Nginx (Step 5) is the reverse proxy.
NoNewPrivileges=true
PrivateTmp=yes
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/var/lib/jabali-filebrowser /home /run/jabali-filebrowser
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=~@clock @debug @module @mount @obsolete @privileged @raw-io @reboot @swap

[Install]
# NOTE: Do NOT enable this service yet. Step 4/5 will wire up the SSO handoff
# and nginx reverse-proxy before we enable it. For now, just make the unit
# available; the service will be started later via systemctl enable/start.
WantedBy=multi-user.target
