# install/systemd/jabali-spam-rules-update.timer
[Unit]
Description=Weekly refresh of Stalwart spam-filter rules (Sun 03:15 host-local)
Documentation=https://github.com/stalwartlabs/spam-filter

[Timer]
# Sunday 03:15 host-local. Off-peak; ~1h before the 04:30 backup retention
# sweep, ~15min after the 03:00 maldet scans. spam-filter rules change at
# upstream's cadence (irregular, days-to-weeks); weekly is generous enough
# that we're never more than 7 days stale on heuristic updates.
OnCalendar=Sun *-*-* 03:15:00
# Spread fires by up to 30min so a hosting fleet doesn't hammer github
# all at the same minute. Cheap good-citizen knob.
RandomizedDelaySec=30min
# Fire missed runs after a host wakes from sleep / suspend / power-cut.
Persistent=true
Unit=jabali-spam-rules-update.service

[Install]
WantedBy=timers.target
