# install/systemd/jabali.slice
[Unit]
Description=Jabali Panel root slice
Before=slices.target

[Slice]
CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
