# install/systemd/jabali-stalwart.service
[Unit]
Description=Jabali Mail (Stalwart SMTP + IMAP + JMAP)
Documentation=https://stalw.art/
After=network-online.target mariadb.service
Wants=network-online.target
PartOf=jabali.slice

[Service]
Type=simple
User=jabali-mail
Group=jabali-mail
# Supplementary group so Stalwart can read /etc/jabali-panel/stalwart-admin.token
# (owned jabali:jabali-mail, mode 0640) and /etc/jabali-panel/dkim/*.key
# (owned jabali:jabali, mode 0600 — Stalwart reads per-domain DKIM keys via
# its config loader).
SupplementaryGroups=jabali
Slice=jabali.slice

ExecStart=/usr/local/bin/stalwart --config=/etc/stalwart/config.toml

Restart=on-failure
RestartSec=5s
WatchdogSec=30s

# Security hardening (matches jabali-kratos.service shape).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/stalwart /etc/stalwart /etc/jabali-panel /run/jabali

# Stalwart binds to ports 25, 465, 587, 993 (external mail protocols) and
# 8446 (loopback JMAP). Only the loopback binding requires no capability;
# binding to the mail ports as non-root needs CAP_NET_BIND_SERVICE.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

StandardOutput=journal
StandardError=journal

[Install]
# Intentionally NOT enabled here — install.sh leaves the unit disabled
# until the first panel `domain.email_enable` call (from the agent).
# See plan §1 and ADR-0041.
WantedBy=multi-user.target
